# etc/profile.d/60-io.neovim.neovim-nightly.sh
# NOTE:
# We have configurate environment variables in /etc/environment.d.
# And we need to export PATH here again
# as PATH is the only one which is overwrited by /etc/profile.

export PATH
PATH=/opt/io.neovim.neovim-nightly/bin:$PATH
