# etc/profile.d/60-io.neovim.neovim-nightly.sh
if [ -f /etc/environment.d/60-io.neovim.neovim-nightly.conf ]; then
	source /etc/environment.d/60-io.neovim.neovim-nightly.conf
fi
export PATH
export MANPATH
export XDG_DATA_DIRS
